# 服务器/rime-server.service.example
[Unit]
Description=Rime Configuration Sync Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=gk-pi
Group=gk-pi
WorkingDirectory=/home/gk-pi/rime_server
ExecStart=/home/gk-pi/rime_server/start_server.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rime-server
Environment=PYTHONUNBUFFERED=1

# 安全设置
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/home/gk-pi/rime_server/logs /home/gk-pi/rime_server/sync /home/gk-pi/rime_server/runtime /home/gk-pi/rime_server/backups

[Install]
WantedBy=multi-user.target